# src/assistant/crontab
* * * * * /usr/local/bin/python /app/ingest.py >> /var/log/cron.log 2>&1
