# src/assistant/crontab
*/15 * * * * /app/scripts/ingest.sh >>/var/log/cron.log 2>&1
